# flarum/rootfs/etc/crontabs/root
* * * * * cd /var/www; su-exec www-data:www-data flarum schedule:run >/dev/stdout 2>/dev/stderr
